# ops/com.sahjony.promotion-evidence-producer.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.sahjony.promotion-evidence-producer</string>
  <key>WorkingDirectory</key><string>/Users/REPLACE_ME/SAHJONY-AI-TRADING</string>
  <key>ProgramArguments</key><array>
    <string>/Users/REPLACE_ME/SAHJONY-AI-TRADING/.venv/bin/python</string>
    <string>scripts/produce_operational_evidence.py</string>
  </array>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/REPLACE_ME/SAHJONY-AI-TRADING/logs/promotion-producer.out.log</string>
  <key>StandardErrorPath</key><string>/Users/REPLACE_ME/SAHJONY-AI-TRADING/logs/promotion-producer.err.log</string>
</dict></plist>
